# torque_adopt.socket
[Unit]
Description=Torque Adopt Socket

[Socket]
MaxConnections=4096
ListenStream=411
Accept=yes

[Install]
WantedBy=sockets.target

